# pbuilder/templates/client/cron.daily/pbuilder.erb
#!/bin/sh

/usr/local/bin/__pbuilder.sh update
